# deploy/bugglo.service
[Unit]
Description=Bugglo — agentic AI for Robinhood Chain
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bugglo
WorkingDirectory=/srv/bugglo
EnvironmentFile=/srv/bugglo/.env.local
Environment=NODE_ENV=production

# Port 3100, NOT 3000 — SpacetimeDB already owns 3000 on this box.
#
# ONE process. Not a cluster.
#
# The spend guard in lib/rateLimit.js keeps its counters in this process's memory, so the
# caps are per-process by construction: run four workers and ENGINE_GLOBAL_USD_PER_DAY=5
# quietly becomes a $20/day ceiling — the exact bill the cap exists to prevent. Move the
# guard to Redis BEFORE ever adding a second worker.
ExecStart=/usr/bin/npm run start -- -p 3100

# The agent spawns MCP servers (robinx-mcp, boo-crypto, dexscreener) as child processes.
# control-group reaping stops a restart from orphaning them holding stdio pipes.
KillMode=control-group
Restart=always
RestartSec=5

# Those MCP children are third-party npm code, fetched at runtime. Keep them away from the
# rest of the box — the Mythreach game servers and SpacetimeDB live here too.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelTunables=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bugglo

[Install]
WantedBy=multi-user.target
