# deploy/bugglo.service
[Unit]
Description=Bugglo — agentic AI for Robinhood Chain
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bugglo
WorkingDirectory=/srv/bugglo
EnvironmentFile=/srv/bugglo/.env.local
Environment=NODE_ENV=production
Environment=PORT=3000

# ONE process. Not a cluster, not N replicas.
#
# The spend guard in lib/rateLimit.js keeps its counters in this process's memory, so the
# caps are per-process by construction. Run four workers and ENGINE_GLOBAL_USD_PER_DAY=5
# silently becomes a $20/day ceiling — the exact bill you set the cap to prevent. If you
# ever do need to scale out, move the guard to Redis FIRST, then add workers.
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000

# The agent spawns MCP servers (robinx-mcp, boo-crypto, dexscreener) as child processes and
# they must be reaped with it, or a restart leaves orphans holding stdio pipes.
KillMode=control-group
Restart=always
RestartSec=5

# The MCP children are untrusted code from npm. Give them as little of the host as possible.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/bugglo/.next /srv/bugglo/node_modules /home/bugglo/.npm

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bugglo

[Install]
WantedBy=multi-user.target
